# docker/config/supervisord.conf
[supervisord]
nodaemon=true

[program:emqttd]
command=/mnt/emqttd/bin/eqmttd start
